# TT-7C0/Release_AAVK/EMG2926AAVK-luci-0.9+svn6252/i18n/italian/root/etc/uci-defaults/luci-i18n-italian
#!/bin/sh
. /etc/functions.sh
include /lib/config

uci_set luci languages it Italian
uci_commit luci
