# TT-7C0/Release_AAVK/package/base-files/platform-dependency-files/EMG2926AAVK/etc/hotplug.d/iface/13-RA_status
#!/bin/sh

. /etc/functions.sh
include /lib/config

defaultwan=$(uci get network.general.defaultWan)

ipv6Enable=$(uci_get network ${INTERFACE} ipv6)

[ "$ipv6Enable" == "1" ] && [ "$ACTION" = "ifup" -a "${INTERFACE:0:3}" = "$defaultwan" ] && /etc/init.d/RA_status restart

