# TT-7C0/Release_AAVK/package/base-files/platform-dependency-files/EMG2926AAVK/etc/hotplug.d/iface/61-qos
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
#!/bin/sh

if [ "$ACTION" = "ifup" ] ; then	

	##reload qos

	/sbin/configure_qos restart


fi
