# TT-7C0/Release_AAVK/EMG2926AAVK-luci-0.9+svn6252/themes/zyxel/root/etc/uci-defaults/luci-theme-zyxel
#!/bin/sh
. /etc/functions.sh
include /lib/config

uci_set luci themes ZyXEL "/luci-static/zyxel"
uci_commit luci
