# TT-7C0/Release_AAVK/package/base-files/platform-dependency-files/EMG2926AAVK/etc/hotplug.d/block/09-usb-app-stop
#!/bin/sh 
# Copyright (C) 2009 OpenWrt.org

. /bin/define-usb-slot "$DEVPATH"

case "$ACTION" in
	add)	
		#/bin/usb-app "start"

	;;
	remove)	
		if df | grep "${USB_SLOT}_${USB_NODE}"; then
			/bin/usb-app "stop"
		fi
	;;
esac	

