# TT-7C0/Release_AAVK/package/network/utils/iptables/files/portTrigger.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2008-2010 OpenWrt.org

START=97

start() {
	### applying after WAN connection up ###
	exit 0
}

stop() {
	portTrigger stop
}

restart() {
	portTrigger restart
}

reload() {
	portTrigger reload
}
