# TT-7C0/Release_AAVK/EMG2926AAVK-luci-0.9+svn6252/i18n/german/root/etc/uci-defaults/luci-i18n-german
#!/bin/sh

. /etc/functions.sh
include /lib/config

uci_set luci languages de Deutsch
uci_commit luci
