# TT-7C0/Release_AAVK/package/base-files/platform-dependency-files/EMG2926AAVK/etc/init.d/usb
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org 

START=39
start() {
	[ -d /proc/bus/usb ] && {
		/bin/mount -t usbfs none /proc/bus/usb
	}

	usb_ctrl 0 disable
	sleep 1
	usb_ctrl 0 enable
	
}
