# TT-7C0/Release_AAVK/EMG2926AAVK-luci-0.9+svn6252/i18n/chinese_traditional/root/etc/uci-defaults/luci-i18n-chinese_traditional
#!/bin/sh
. /etc/functions.sh
include /lib/config

uci_set luci languages zh Chinese
uci_commit luci
