# TT-7C0/Release_AAVK/EMG2926AAVK-luci-0.9+svn6252/i18n/spanish/root/etc/uci-defaults/luci-i18n-spanish
#!/bin/sh
. /etc/functions.sh
include /lib/config

uci_set luci languages es Spanish
uci_commit luci
