# TT-7C0/Release_AAVK/package/base-files/files/etc/hotplug.d/iface/10-dns
##check mode
local sys_mode=$(get_sys_mode)
[ "$sys_mode" == "2" -o "$sys_mode" == "3" ] && {

if [ "$ACTION" = "ifup" -a "$INTERFACE" = "lan" ]; then
	/sbin/update_sys_dns_apmode
fi

} || {

if [ "$ACTION" = "ifup" -a "$INTERFACE" = "wan" ]; then
	/sbin/update_sys_dns
	/sbin/update_lan_dns
	/etc/init.d/dnsmasq stop
	/etc/init.d/dnsmasq start
fi

}
