# TT-7C0/Release_AAVK/EMG2926AAVK-luci-0.9+svn6252/i18n/french/root/etc/uci-defaults/luci-i18n-french
#!/bin/sh
. /etc/functions.sh
include /lib/config

uci_set luci languages fr Francais
uci_commit luci
