# TT-7C0/Release_AAVK/package/base-files/files/etc/inittab
::sysinit:/etc/init.d/rcS S boot
::shutdown:/etc/init.d/rcS K shutdown
#ttyS0::askfirst:/bin/ash --login
tty1::askfirst:/bin/ash --login
ttyS0::askfirst:/sbin/getty -L ttyS0 115200 vt100
